# supervisord.conf
[supervisord]
nodaemon=false

[program:mongodb]
command=/bin/bash -c "exec service mongodb start"

[program:server]
command=/bin/bash -c "exec ./main"


