# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/bin/bash -c "exec ./main"