# supervisord.conf
[supervisord]
nodaemon=false

[program:mongodb]
command=/bin/bash -c "exec service mongod start"

[program:server]
command=/bin/bash -c "exec ./main"


