# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=/bin/bash -c "exec ./main"


#[program:mongodb]
#command=/bin/bash -c "exec service mongod start"